# utils/dockerd/files/etc/hotplug.d/iface/dockerd

[ "$ACTION" = "ifup" -a -n "$DEVICE" ] || exit 0

/etc/init.d/dockerd running && /etc/init.d/dockerd reload
